# systemd/gluster-blockd.service.in
[Unit]
Description=Gluster block storage utility
Requisite=glusterd.service
Requires=rpcbind.service
BindsTo=gluster-block-target.service
Wants=targetclid.service
After=gluster-block-target.service rpcbind.service

[Service]
Type=simple
Environment="GB_GLFS_LRU_COUNT=5"
Environment="GB_LOG_LEVEL=INFO"
EnvironmentFile=-@sysconfigdir@/gluster-blockd
ExecStart=@prefix@/sbin/gluster-blockd --glfs-lru-count $GB_GLFS_LRU_COUNT --log-level $GB_LOG_LEVEL $GB_EXTRA_ARGS
KillMode=process
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target
